# debian/munin.upstart
# vim:set ft=upstart ts=2 et:
description "Munin"
author "James Page <james.page@ubuntu.com>"

start on runlevel [2345]
stop on runlevel [!2345]

task

script
  # Ensure /var/run/munin exists with correct permissions
  install -m 0755 -o munin -g munin -d /var/run/munin
end script
