# debian/munin-async.upstart
# vim: set ft=upstart ts=2 et:
description "Munin asynchronous server"
author "James Page <james.page@ubuntu.com>"

start on started munin-node
stop on stopping munin-node

respawn

script
  [ -r /etc/default/munin-asyncd ] && . /etc/default/munin-asyncd
  exec start-stop-daemon --start --chuid munin-async \
    --exec /usr/share/munin/munin-asyncd -- $DAEMON_ARGS
end script
