# debian/munin-node.upstart
# vim: set ft=upstart ts=2 et:
description "munin-node"
author "Chuck Short <zulcss@ubuntu.com>"

start on (filesystem and net-device-up IFACE=lo)
stop on runlevel [!2345]

expect fork
respawn

pre-start script
  install -m 0755 -o munin -g munin -d /var/run/munin
end script

script
  [ -r /etc/default/munin-node ] && . /etc/default/munin-node
  exec /usr/sbin/munin-node $DAEMON_ARGS
end script
