# dists/freebsd/munin-node/pkg-plist
etc/munin/munin-node.conf
etc/munin/plugin-conf.d/plugins.conf
etc/rc.d/munin-node.sh
lib/plugins/apache_accesses
lib/plugins/apache_processes
lib/plugins/apache_volume
lib/plugins/bind9
lib/plugins/cpu
lib/plugins/df
lib/plugins/df_inode
lib/plugins/exim_mailqueue
lib/plugins/exim_mailstats
lib/plugins/hddtemp
lib/plugins/if_
lib/plugins/if_errcoll_
lib/plugins/ircu
lib/plugins/load
lib/plugins/mailman
lib/plugins/memory
lib/plugins/munin_graph
lib/plugins/munin_update
lib/plugins/mysql_bytes
lib/plugins/mysql_isam_space_
lib/plugins/mysql_queries
lib/plugins/mysql_slowqueries
lib/plugins/mysql_threads
lib/plugins/named
lib/plugins/netstat
lib/plugins/open_files
lib/plugins/plugins.history
lib/plugins/pop_stats
lib/plugins/postfix_mailqueue
lib/plugins/postfix_mailstats
lib/plugins/postfix_mailvolume
lib/plugins/processes
lib/plugins/ps_
lib/plugins/psu_
lib/plugins/samba
lib/plugins/snmp__if_
lib/plugins/spamstats
lib/plugins/squid_cache
lib/plugins/squid_icp
lib/plugins/squid_requests
lib/plugins/squid_traffic
lib/plugins/surfboard
lib/plugins/swap
lib/plugins/sybase_space
lib/plugins/vmstat
sbin/munin-node
sbin/munin-node-configure
sbin/munin-run
@unexec rmdir %D/var/munin/plugin-state
@unexec rmdir %D/var/munin
@dirrm lib/plugins
@unexec rmdir %D/etc/munin/plugins
@unexec rmdir %D/etc/munin/plugin-conf.d
@unexec rmdir %D/etc/munin
