# dists/freebsd/munin-node/pkg-plist
@comment etc/munin/VERSION.node must remain for upgrade version check
etc/munin/munin-node.conf
etc/munin/plugin-conf.d/plugins.conf
etc/rc.d/munin-node.sh
sbin/munin-node
sbin/munin-node-configure
sbin/munin-run
%%DATADIR%%/plugins/acpi
%%DATADIR%%/plugins/amavis
%%DATADIR%%/plugins/apache_accesses
%%DATADIR%%/plugins/apache_processes
%%DATADIR%%/plugins/apache_volume
%%DATADIR%%/plugins/apc_envunit_
%%DATADIR%%/plugins/bind9
%%DATADIR%%/plugins/cpu
%%DATADIR%%/plugins/df
%%DATADIR%%/plugins/df_inode
%%DATADIR%%/plugins/exim_mailqueue
%%DATADIR%%/plugins/exim_mailstats
%%DATADIR%%/plugins/folding@home
%%DATADIR%%/plugins/hddtemp
%%DATADIR%%/plugins/hddtemp2
%%DATADIR%%/plugins/if_
%%DATADIR%%/plugins/if_errcoll_
%%DATADIR%%/plugins/ipac-ng
%%DATADIR%%/plugins/ircu
%%DATADIR%%/plugins/load
%%DATADIR%%/plugins/mailman
%%DATADIR%%/plugins/memory
%%DATADIR%%/plugins/multips
%%DATADIR%%/plugins/munin_graph
%%DATADIR%%/plugins/munin_update
%%DATADIR%%/plugins/mysql_bytes
%%DATADIR%%/plugins/mysql_isam_space_
%%DATADIR%%/plugins/mysql_queries
%%DATADIR%%/plugins/mysql_slowqueries
%%DATADIR%%/plugins/mysql_threads
%%DATADIR%%/plugins/named
%%DATADIR%%/plugins/netstat
%%DATADIR%%/plugins/open_files
%%DATADIR%%/plugins/plugins.history
%%DATADIR%%/plugins/pm3users
%%DATADIR%%/plugins/pop_stats
%%DATADIR%%/plugins/postfix_mailqueue
%%DATADIR%%/plugins/postfix_mailstats
%%DATADIR%%/plugins/postfix_mailvolume
%%DATADIR%%/plugins/processes
%%DATADIR%%/plugins/ps_
%%DATADIR%%/plugins/psu_
%%DATADIR%%/plugins/samba
%%DATADIR%%/plugins/sendmail_mailqueue
%%DATADIR%%/plugins/sendmail_mailstats
%%DATADIR%%/plugins/sendmail_mailtraffic
%%DATADIR%%/plugins/snmp__df
%%DATADIR%%/plugins/snmp__fc_if_
%%DATADIR%%/plugins/snmp__fc_if_err_
%%DATADIR%%/plugins/snmp__if_
%%DATADIR%%/plugins/snmp__if_err_
%%DATADIR%%/plugins/snmp__load
%%DATADIR%%/plugins/snmp__processes
%%DATADIR%%/plugins/snmp__users
%%DATADIR%%/plugins/spamstats
%%DATADIR%%/plugins/squid_cache
%%DATADIR%%/plugins/squid_icp
%%DATADIR%%/plugins/squid_requests
%%DATADIR%%/plugins/squid_traffic
%%DATADIR%%/plugins/surfboard
%%DATADIR%%/plugins/swap
%%DATADIR%%/plugins/sybase_space
%%DATADIR%%/plugins/vmstat
@unexec rmdir %D/var/munin/plugin-state 2>/dev/null || true
@unexec rmdir %D/var/munin 2>/dev/null || true
@unexec rmdir %D/%%DATADIR%% 2>/dev/null || true
@unexec rmdir %D/etc/munin/plugin-conf.d 2>/dev/null || true
@unexec rmdir %D/etc/munin 2>/dev/null || true
@unexec echo "If you want to *completely* remove munin-node, you have to"
@unexec echo "    rm -rf %D/var/munin and %D/etc/munin %D/www/data/munin"
@unexec echo "They have to be left on the system because you may be performing"
@unexec echo "an update. This process can't tell."
