# dists/freebsd/munin-node/files/munin-node.in
#!/bin/sh
#
# $FreeBSD: head/sysutils/munin-node/files/munin-node.in 447210 2017-08-03 13:34:57Z mat $
#

# PROVIDE: munin-node
# REQUIRE: DAEMON
# BEFORE:  cron
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable munin-node:
# munin_node_enable (bool):      Set to "NO" by default.
#                                Set it to "YES" to enable munin-node
# munin_node_config (path):      Set to "%%PREFIX%%/etc/munin/munin-node.conf" by default.
#

. /etc/rc.subr

name=munin_node
desc="munin node daemon"
rcvar=munin_node_enable

load_rc_config $name

command="%%PREFIX%%/sbin/munin-node"
command_interpreter="%%LOCALBASE%%/bin/perl"
start_precmd=find_pidfile
status_precmd=find_pidfile
stop_precmd=find_pidfile

munin_node_enable=${munin_node_enable:-NO}
munin_node_config=${munin_node_config:-"%%PREFIX%%/etc/munin/munin-node.conf"}

find_pidfile()
{
	get_pidfile_from_conf pid_file $munin_node_config
	pidfile="$_pidfile_from_conf"
}

run_rc_command $1
